# mitmproxy/root/etc/services.d/mitmproxy/run
#!/usr/bin/with-contenv bashio
set -e
OPTIONS=$(bashio::jq "${__BASHIO_ADDON_CONFIG}" '.options | map("\(.name)=\(.value|tostring)") | [.[] | "--set "  + .] |  join(" ")')
HOST=$(bashio::host.hostname)

bashio::log.info "Start mitmproxy with options: $OPTIONS"
bashio::log.info "For the onboarding page visit http://${HOST}:8082"
# shellcheck disable=SC2086
mitmweb $OPTIONS \
  --set 'web_host=0.0.0.0' \
  --set 'web_port=9090' \
  --set 'listen_port=8080' \
  --set 'confdir=/data/mitmproxy' \
  --set 'onboarding_port=8082' \
  --set "onboarding_host=$HOST" \
  --set 'web_open_browser=false'
  --no-web-open-browser
