# pixelflut/root/etc/s6-overlay/s6-rc.d/shoreline/run
#!/command/with-contenv bashio
set -e

width=$(bashio::config 'width')
height=$(bashio::config 'height')
description=$(bashio::config 'description')

bashio::log.info 'Start shoreline pixelflut server'
shoreline -f "vnc,port=5901" -f "statistics,listen=0.0.0.0,port=1235" -b 0.0.0.0 -p 1234 -t '/arial.ttf' -w "$width" -h "$height" -d "$description" 
