# pigpio/root/etc/services.d/pigpio/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Started pigpio'
pigpiod -g -f
