# valetudo-mapper/root/etc/s6-overlay/s6-rc.d/valetudo-mapper/run
#!/command/with-contenv bashio
set -e

# Enable Jemalloc for better memory handling
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

bashio::log.info 'Setup config'
jq -f /app/config-template.jq /data/options.json > /app/config.json

bashio::log.info 'Start valetudo-mapper'
cd /app
npm start
