# pixelflut/root/etc/s6-overlay/s6-rc.d/noVNC/run
#!/command/with-contenv bashio
set -e
bashio::net.wait_for 5900

bashio::log.info 'Start noVNC'
bashio::log.info "Binding to Pixelflut canvas"
/noVNC/utils/novnc_proxy --listen 6080 --vnc localhost:5900
