# valetudo-mapper/root/etc/services.d/valetudo-mapper/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Start valetudo-mapper'
cd /app
npm start
