# valetudo-mapper/root/etc/services.d/valetudo-mapper/run
#!/usr/bin/with-contenv bashio
set -e

# Enable Jemalloc for better memory handling
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

bashio::log.info 'Start valetudo-mapper'
cd /app
npm start
