# pigpio/root/etc/s6-overlay/s6-rc.d/pigpio/run
#!/command/with-contenv bashio
set -e

additional_options=$(bashio::config 'additional_options')

bashio::log.info 'Started pigpio'
# shellcheck disable=SC2086
pigpiod -g -f $additional_options
