# container-stats/root/etc/services.d/stats/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Start docker statistics'

export STATS_UPDATE_INTERVAL=$(bashio::config 'update_interval')

cd /container-stats
node stats.js
