# container-stats/root/etc/services.d/stats/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Start docker statistics'

if bashio::config.has_value 'update_interval'; then
	export STATS_UPDATE_INTERVAL=$(bashio::config 'update_interval')
fi

if [ ! -d /container-stats/db ] && bashio::config 'persistent'; then
	[ ! -d /data/db ] && mkdir -p /data/db
	ln -s /data/db/ /container-stats/db
else
	mkdir -p /container-stats/db
fi

cd /container-stats
node stats.js
