# MaryTTS/root/etc/services.d/marytts/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start MaryTTS (taking some time, will output the port when up)'
cd /app
bin/marytts-server
