# ICantBelieveItsNotValetudo/root/etc/s6-overlay/s6-rc.d/Server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info 'Start Proxywebserver'
cd /server
npm start
