# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing-setup/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

# determine IP address for api requests
ip=$(ifconfig | grep "inet " | grep -Fv 127.0.0.1 | awk '{sub("addr:",""); print $2}' | head -n 1)

if [ -z "$ip" ]
then
  ip=0.0.0.0
  bashio::log.info "Couldn't detect local syncthing address. Falling back to: $ip"
else
  bashio::log.info "Detected syncthing address: $ip"
fi

bashio::net.wait_for 8384 "$ip"
bashio::log.info 'Post-Start syncthing setup'

apikey=$(grep -o '<apikey>[^<]*' /config/config.xml | sed 's/<apikey>//;s/<\/apikey>//')
syncthing cli --gui-address="$ip:8384" --gui-apikey="$apikey" config gui insecure-admin-access set true
