# syncthing/root/etc/services.d/syncthing/run
#!/usr/bin/with-contenv bashio
set -e

# Enable Jemalloc for better memory handling
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

export STNOUPGRADE=1

bashio::log.info 'Start syncthing'
syncthing -no-browser -no-restart -home=/data/config -gui-address=0.0.0.0:8384
