# container-stats/root/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Start webserver'
cd /container-stats
node httpd.js
