# ympd/root/etc/s6-overlay/s6-rc.d/ympd/run
#!/command/with-contenv bashio
set -e

host=$(bashio::config 'host')
port=$(bashio::config 'port')

bashio::log.info 'Start ympd'
ympd -w 6680 -h "$host" -p "$port"
