# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info 'Intializing Syncthing'

# determine IP address to serve Syncthing's GUI on
ip=$(ifconfig | grep "inet " | grep -Fv 127.0.0.1 | awk '{sub("addr:",""); print $2}' | head -n 1)

if [ -z "$ip" ]
then
  ip=0.0.0.0
  bashio::log.info "Couldn't detect IP address. Falling back to: $ip"
else
  bashio::log.info "Detected IP address: $ip"
fi

# move config and DB state from old to new locations if necessary, cf. https://github.com/Poeschl/Hassio-Addons/pull/450
if [ -d /data/config ]
then
  for CONF_FILE in config.xml cert.pem key.pem https-cert.pem https-key.pem
  do
    [ -f "/data/config/${CONF_FILE}" ] && mv -fn "/data/config/${CONF_FILE}" /config/ && rm -f "/data/config/${CONF_FILE}"
  done
  mv -fn /data/config/* /data/ && rm -rf /data/config
fi

bashio::log.info 'Starting Syncthing'
syncthing --gui-address="$ip:8384" --no-browser
