# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

export STNOUPGRADE=1

bashio::log.info 'Setup config'
mkdir -p /data/config

bashio::log.info 'Start syncthing'
syncthing --no-browser --no-restart --home=/data/config --gui-address=0.0.0.0:8384
