# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info 'Intializing Syncthing'

ip=$(ifconfig | grep "inet " | grep -Fv 127.0.0.1 | awk '{sub("addr:",""); print $2}' | head -n 1)

if [ -z "$ip" ]
then
  ip=0.0.0.0
  bashio::log.info "Couldn't detect IP address. Falling back to: $ip"
else
  bashio::log.info "Detected IP address: $ip"
fi

bashio::log.info 'Starting Syncthing'
syncthing --gui-address="$ip:8384" --no-browser
