# vnc-viewer/root/etc/services.d/noVNC/run
#!/usr/bin/with-contenv bashio
set -e

if bashio::config.has_value 'vnc_server_password'; then
    bashio::log.info 'Inject vnc connection password'
    password=$(bashio::config 'vnc_server_password')
    sed -i "s/WebUtil.getConfigVar('password');/'${password}';/" /noVNC/app/ui.js
fi

server_host=$(bashio::config 'vnc_server_host')
server_port=5900
if bashio::config.has_value 'vnc_server_port'; then
    server_port=$(bashio::config 'vnc_server_port')
fi

server="${server_host}:${server_port}"

bashio::log.info 'Start noVNC'
bashio::log.info "Connecting to VNC $server"
/noVNC/utils/launch.sh --listen 6080 --vnc $server
