# picoTTS/root/etc/services.d/picotts/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start PicoTTS webserver'
cd /app
python3 -u src/picowrapper.py
